# default-data/config/supervisord.conf
; supervisor config file

[unix_http_server]
file=/container/system/run/supervisor.sock  	; (the path to the socket file)
chmod=0766                      	; socket file mode (default 0700)

[supervisord]
logfile=/container/logs/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid 			; (supervisord pidfile;default supervisord.pid)
childlogdir=/container/logs/supervisor             ; ('AUTO' child log dir, default $TEMP)
umask=022

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///container/system/run/supervisor.sock	; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[supervisord]
nodaemon=true

[program:rsyslog]
command=rsyslogd -n -f /container/config/rsyslog.conf
stdout_logfile=/container/logs/supervisor/%(program_name)s-stdout.log
stderr_logfile=/container/logs/supervisor/%(program_name)s-stderr.log

[program:nginx]
command=/usr/sbin/nginx
stdout_logfile=/container/logs/supervisor/%(program_name)s-stdout.log
stderr_logfile=/container/logs/supervisor/%(program_name)s-stderr.log

[program:cron]
command=/bin/bash -c "touch /etc/crontab && cron -f"
stdout_logfile=/container/logs/supervisor/%(program_name)s-stdout.log
stderr_logfile=/container/logs/supervisor/%(program_name)s-stderr.log

[program:php]
command=/sbin/php-fpm --nodaemonize --fpm-config /container/config/php/php-fpm.conf
stdout_logfile=/container/logs/supervisor/%(program_name)s-stdout.log
stderr_logfile=/container/logs/supervisor/%(program_name)s-stderr.log

[program:cleartmp]
command=/bin/bash -c "/usr/bin/rm -rf /tmp/*"
startretries=0
startsecs=0
priority=0
autorestart=false

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /container/application/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=8
redirect_stderr=true
stdout_logfile=/container/logs/supervisor/%(program_name)s-worker.log
stopwaitsecs=3600
